# public/Library/LaunchAgents/teohm.sysupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>teohm.sysupdate</string>
    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-i</string>
      <string>-c</string>
      <string>$SHELL --login -c "$HOME/bin/sysupdate"</string>
    </array>
    <key>StandardErrorPath</key>
    <string>/Users/huiming/log/sysupdate.log</string>
    <key>StandardOutPath</key>
    <string>/Users/huiming/log/sysupdate.log</string>
    <key>StartCalendarInterval</key>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </dict>
</plist>
